# website_backend_laravel/.deploy/config/supervisor.conf
[supervisord]
logfile=/dev/null
logfile_maxbytes=0
logfile_backups=0
loglevel=info
nodaemon=true

[program:php-fpm]
command=php-fpm --nodaemonize
autorestart=true
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr

[program:caddy]
command=caddy run --config %(ENV_LARAVEL_PATH)s/.deploy/config/Caddyfile
autorestart=true
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr

[program:cron]
command=crond -l 2 -f
autorestart=true


[program:queue-listen]
command=php %(ENV_LARAVEL_PATH)s/artisan queue:listen --queue=high,medium,admin,low,default --timeout 3600
numprocs=1
autorestart=true
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr


[program:queue-processor0]
command=php %(ENV_LARAVEL_PATH)s/artisan queue:listen --queue=processor-0  --timeout 3600
numprocs=1
autorestart=true
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr

[program:queue-processor1]
command=php %(ENV_LARAVEL_PATH)s/artisan queue:listen --queue=processor-1  --timeout 3600
numprocs=1
autorestart=true
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr


[program:queue-processor2]
command=php %(ENV_LARAVEL_PATH)s/artisan queue:listen --queue=processor-2  --timeout 3600
numprocs=1
autorestart=true
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr


[program:queue-processor3]
command=php %(ENV_LARAVEL_PATH)s/artisan queue:listen --queue=processor-3  --timeout 3600
numprocs=1
autorestart=true
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr

[program:queue-processor4]
command=php %(ENV_LARAVEL_PATH)s/artisan queue:listen --queue=processor-4  --timeout 3600
numprocs=1
autorestart=true
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr

[program:queue-processor5]
command=php %(ENV_LARAVEL_PATH)s/artisan queue:listen --queue=processor-5  --timeout 3600
numprocs=1
autorestart=true
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr

[program:queue-processor6]
command=php %(ENV_LARAVEL_PATH)s/artisan queue:listen --queue=processor-6  --timeout 3600
numprocs=1
autorestart=true
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr

[program:queue-processor7]
command=php %(ENV_LARAVEL_PATH)s/artisan queue:listen --queue=processor-7  --timeout 3600
numprocs=1
autorestart=true
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr

[program:queue-processor8]
command=php %(ENV_LARAVEL_PATH)s/artisan queue:listen --queue=processor-8  --timeout 3600
numprocs=1
autorestart=true
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr

[program:queue-processor9]
command=php %(ENV_LARAVEL_PATH)s/artisan queue:listen --queue=processor-9  --timeout 3600
numprocs=1
autorestart=true
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr

; [program:horizon]
; command=php %(ENV_LARAVEL_PATH)s/artisan horizon
; autorestart=true
; stdout_events_enabled=true
; stderr_events_enabled=true
; stdout_logfile_maxbytes=0
; stderr_logfile_maxbytes=0
; stdout_logfile=/dev/stdout
; stderr_logfile=/dev/stderr
; stopwaitsecs=3600

; [program:websockets]
; command=php %(ENV_LARAVEL_PATH)s/artisan websockets:serve --port=6001
; numprocs=1
; autorestart=true
; stdout_events_enabled=true
; stderr_events_enabled=true
; stdout_logfile_maxbytes=0
; stderr_logfile_maxbytes=0
; stdout_logfile=/dev/stdout
; stderr_logfile=/dev/stderr
